# caddy/rootfs/etc/s6-overlay/s6-rc.d/caddy/run
#!/command/with-contenv bashio
# shellcheck shell=bash

exec /usr/bin/run-caddy.sh

